# deploy/systemd/zeroone.service
[Unit]
Description=Zeroone Control Plane
After=network-online.target xray.service openvpn@company.service openvpn@client-tun1.service
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=main
EnvironmentFile=-/etc/default/zeroone
ExecStart=/usr/local/bin/zeroone -config /usr/local/etc/zeroone/stack.json $ZEROONE_FLAGS
Restart=always
RestartSec=5
WatchdogSec=30s
User=root
LimitNOFILE=1048576
TasksMax=infinity
MemoryMax=512M
MemoryHigh=384M

[Install]
WantedBy=multi-user.target
