# rootfs/etc/systemd/system/vpn-monitor.service
[Unit]
Description=VPN Stack Monitor
After=openvpn@company.service xray.service
Wants=openvpn@company.service xray.service

[Service]
ExecStart=/usr/local/bin/vpn-monitor.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
